# files/ld.so.preload
/usr/lib/libhardened_malloc.so/libhardened_malloc.so
